# conf.d/sms.conf
[program:gunicorn]
command=/home/<project_folder>/sms/bin/start_gunicorn.sh
user=<your_user>
process_name=%(program_name)s
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
